# debian.default/puavo-accounts.upstart
description "puavo-accounts"

start on runlevel [23]
stop on shutdown

pre-start script
    mkdir -p /run/puavo-accounts
    chown puavo:puavo /run/puavo-accounts
end script

chdir /usr/lib/puavo-accounts
env RACK_ENV=production

exec sudo -E -u puavo ruby2.0 /usr/bin/bundle exec puma --workers 4 --port 9494

respawn
